# files/common/base_linux/bin/puppet-agent-first-run.cron
# Run puppet agent
*/5 * * * * root /usr/local/sbin/puppet-agent-run.sh > /dev/null 2>&1
