# files/common/base_linux/bin/puppet-agent-run.cron
# Run puppet agent
*/20 * * * * root /usr/local/sbin/puppet-agent-run.sh > /dev/null 2>&1
